# salt/minion/files/salt-minion-cache-clear.timer.jinja
{%- from "salt/map.jinja" import saltdata %}
[Timer]
OnCalendar={{ saltdata.minion.cachecleaner.oncalendar }}

[Install]
WantedBy=multi-user.target
