# salt/minion/files/salt-minion-cache-clear.timer.jinja
[Timer]
OnCalendar=weekly

[Install]
WantedBy=multi-user.target
